# hl_engine/live-trader.service
[Unit]
Description=HL Engine — APEX Live Trader
Documentation=https://github.com/hollowc2/HyperLiquid
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
User=billy
WorkingDirectory=/opt/Hyperliquid/hl_engine
ExecStartPre=docker compose build --quiet live-trader
ExecStart=docker compose up live-trader
ExecStop=docker compose stop live-trader
Restart=always
RestartSec=15
StartLimitIntervalSec=180
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hl-live-trader

[Install]
WantedBy=multi-user.target
