# hl_engine/live-trader.service
[Unit]
Description=HL Engine — APEX Live Trader
Documentation=https://github.com/hollowc2/hyperliquid
After=network-online.target
StartLimitIntervalSec=180
StartLimitBurst=5

[Service]
Type=simple
User=billy
WorkingDirectory=/opt/hyperliquid/hl_engine
Environment=HL_STATE_DIR=/opt/hyperliquid/hl_engine/data
ExecStart=/opt/hyperliquid/hl_engine/.venv/bin/python -m hl_engine.run_live
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hl-live-trader

[Install]
WantedBy=multi-user.target
