# hl_engine/ts-recorder.service
[Unit]
Description=HL Engine — TimescaleDB Live Recorder
Documentation=https://github.com/hollowc2/hyperliquid
After=docker.service network-online.target
Requires=docker.service
StartLimitIntervalSec=180
StartLimitBurst=5

[Service]
Type=simple
User=billy
WorkingDirectory=/opt/hyperliquid/hl_engine
ExecStartPre=docker compose build --quiet ts-recorder
ExecStart=docker compose up --no-deps ts-recorder
ExecStop=docker compose stop ts-recorder
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hl-ts-recorder

[Install]
WantedBy=multi-user.target
