# hl_engine/ts-recorder.service
[Unit]
Description=HL Engine — TimescaleDB Live Recorder
Documentation=https://github.com/hollowc2/HyperLiquid
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
User=billy
WorkingDirectory=/opt/Hyperliquid/hl_engine
ExecStartPre=docker compose build --quiet ts-recorder
ExecStart=docker compose up ts-recorder
ExecStop=docker compose stop ts-recorder
Restart=always
RestartSec=15
StartLimitIntervalSec=180
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hl-ts-recorder

[Install]
WantedBy=multi-user.target
